# debian/sonic-object-library.postinst
#!/bin/bash



systemctl daemon-reload || /bin/systemctl enable sonic-object-library  || /bin/systemctl start sonic-object-library 
