# debian/sonic-object-library.prerm
#!/bin/sh
systemctl stop sonic-object-library
systemctl disable sonic-object-library
