# debian/sonic-object-library.postinst
#!/bin/sh
systemctl enable sonic-object-library
systemctl start sonic-object-library
